# Modul 3/crontab
* * * * * echo "Hello world at "$(date '+%Y-%m-%d_%H:%M:%S') >> /var/log/cron.log 2>&1
0 */3 * * * /usr/local/bin/backup
#* * * * * /usr/local/bin/backup
